# .config/systemd/user/sway-headless.service
[Unit]
Description=Headless Sway for luth (CC dev workstation)

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=sway
Environment=WLR_BACKENDS=headless
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=WLR_RENDERER=pixman
ExecStart=/usr/bin/sway -c %h/.config/sway/config-headless
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
